# systemd/rfserve_batch.timer
[Unit]
Description=Run Rfserve Batch (daily)

[Timer]
OnCalendar=00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target